# xmastree.openrc
#!/sbin/openrc-run

name="xmastree"
description="PiHut 3D Xmas Tree"
command="/usr/local/bin/xmastree"
command_user="viny"
pidfile="/run/xmastree.pid"
command_background="yes"